# scoreboard.service
[Unit]
Description=Generic Game Scoreboard Backend
Documentation=https://github.com/furkanbaytekin/generic-game-scoreboard-backend
After=network.target

[Service]
Type=simple
User=scoreboard
Group=scoreboard
DynamicUser=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
# Go runtime needs: ProtectControlGroups, ProtectKernelTunables, MemoryDenyWriteExecute disabled.
# SystemCallFilter removed — pure Go runtime (modernc sqlite) uses syscalls
# like rseq, membarrier, futex_waitv not in @default.
WorkingDirectory=/var/lib/scoreboard
StateDirectory=scoreboard
EnvironmentFile=/etc/scoreboard/.env
ExecStart=/usr/local/bin/scoreboard
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
